# meta-ampere/meta-mitchell/recipes-ampere/platform/ampere-misc/power-manager/power-cap-exceeds-limit.service
[Unit]
Description= Notify that total power consumption exceeds the power cap

[Service]
ExecStart=/usr/sbin/power_cap_exceed_limit.sh
Type=oneshot
