# meta-ampere/meta-mitchell/recipes-ampere/platform/ampere-misc/power-manager/power-cap-drops-below-limit.service
[Unit]
Description= Notify that total power consumption drops below the power cap

[Service]
ExecStart=/usr/sbin/power_cap_drop_below_limit.sh
Type=oneshot
