# meta-ampere/meta-common/recipes-phosphor/flash/phosphor-software-manager/allow-reboot-actions.service
[Unit]
Description=Ampere service to allow reboot actions

[Service]
Type=oneshot
ExecStart=/usr/sbin/host_power_control_flashing.sh reboot true