# meta-ampere/meta-common/recipes-phosphor/flash/phosphor-software-manager/prevent-reboot-actions.service
[Unit]
Description=Ampere service to prevent reboot actions

[Service]
Type=oneshot
ExecStart=/usr/sbin/ampere_power_control_lock.sh reboot false
